# crontab
0 * * * * cd ~/opt/mapperserver && python3 draw_map.py
